# files/container-traefik-proxy-8080.socket
[Socket]
ListenStream=192.168.178.39:8080
Service=container-traefik-proxy@8080.service

[Install]
WantedBy=sockets.target