# files/container-traefik-proxy-80.socket
[Socket]
ListenStream=192.168.178.39:80
Service=container-traefik-proxy@80.service

[Install]
WantedBy=sockets.target